# displaylink.target
[Unit]
Description=DisplayLink environment target
Requires=dl-switch.service
Requires=graphical.target
